# brew-base/etc/profile.d/gem.sh
GEM_HOME="/opt/gem/home"
PATH="$GEM_HOME/bin:$PATH"

export GEM_HOME
export PATH
