# brew-base/etc/profile.d/groovy.sh
GROOVY_HOME="/usr/local/opt/groovy/libexec"
export GROOVY_HOME
