# brew-base/etc/profile.d/userbin.sh
USER_BIN="`echo ~/bin`"
PATH="$USER_BIN":"$PATH"

export PATH
